# contrib/librespot.user.service
[Unit]
Description=Librespot

[Service]
Restart=always
RestartSec=12
ExecStart=/usr/bin/librespot --name "%n@%H"

[Install]
WantedBy=default.target
